# dist/etc/systemd/system/myservice.service
[Unit]
Description=My awesome service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mybinary

[Install]
WantedBy=multi-user.target
